# packages/network/wireguard-tools/system.d/wireguard-defaults.service
[Unit]
Description=WireGuard defaults
After=local-fs.target

ConditionPathExists=|!/storage/.cache/wireguard/privatekey
ConditionPathExists=|!/storage/.cache/wireguard/publickey
ConditionPathExists=|!/storage/.cache/wireguard/preshared

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/bin/wireguard-setup'
RemainAfterExit=yes
